# docker/supervisord.conf
[supervisord]
nodaemon=true
user=toolrecall
logfile=/dev/stdout
logfile_maxbytes=0

[program:daemon]
command=toolrecall daemon --foreground
directory=/toolrecall
user=toolrecall
environment=TOOLRECALL_CACHE_DB="/data/cache.db",TOOLRECALL_KNOWLEDGE_DB="/data/knowledge.db",TOOLRECALL_UDS_PATH="/data/tc.sock"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:proxy]
command=toolrecall serve
directory=/toolrecall
user=toolrecall
environment=TOOLRECALL_CACHE_DB="/data/cache.db",TOOLRECALL_KNOWLEDGE_DB="/data/knowledge.db",TOOLRECALL_UDS_PATH="/data/tc.sock"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0