# scripts/toolrecall-daemon.service
[Unit]
Description=ToolRecall Daemon — Tool-Output Cache + MCP Multiplexer
Documentation=https://github.com/whiskybeer/toolrecall
After=network.target

[Service]
Type=simple
ExecStart=python3 -m toolrecall daemon --foreground
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=soft
ProtectSystem=full

[Install]
WantedBy=multi-user.target
