# scripts/toolrecall-daemon.service
[Unit]
Description=ToolRecall Daemon — Tool-Output Cache + MCP Multiplexer
Documentation=https://github.com/whiskybeer/toolrecall
After=network.target

[Service]
Type=simple
# ExecStartPre: Prevent restart loop when daemon is already running.
# Checks PID file (~/.toolrecall/daemon.pid) — if process alive, exits 0
# so systemd skips start. If stale or absent, exits 0 to continue.
# Always exits 0 (never exit 1) because exit 1 + Restart=on-failure
# would trigger a restart loop.
ExecStartPre=/usr/bin/env bash -c 'if [ -f %h/.toolrecall/daemon.pid ]; then pid=$(cat %h/.toolrecall/daemon.pid 2>/dev/null) && kill -0 "$pid" 2>/dev/null && echo "Daemon already running (PID $pid)" && exit 0; fi; exit 0'
ExecStart=python3 -m toolrecall daemon --foreground
PIDFile=%h/.toolrecall/daemon.pid
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=soft
ProtectSystem=full

[Install]
WantedBy=multi-user.target
